# auth-monitor.service
[Unit]
Description=Twitter Bot Authentication Monitor
After=network.target twitter-bot.service
Wants=twitter-bot.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/Twitter-bot
ExecStart=/usr/bin/python3 auth_monitor.py
# DISABLED - Auth monitor causes false positives and deactivates bots
# Restart=always
Restart=no
RestartSec=30
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/root/Twitter-bot
Environment=PYTHONUNBUFFERED=1

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/Twitter-bot

[Install]
WantedBy=multi-user.target
