# webhook-listener.service
[Unit]
Description=GitHub Webhook Listener
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Twitter-bot
ExecStart=/usr/bin/python3 /root/Twitter-bot/webhook_listener.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=webhook-listener

[Install]
WantedBy=multi-user.target